# profile.d/ssh-agent.sh
export SSH_AUTH_SOCK=$HOME/.ssh/ssh-agent.sock
